# docker/scheduler-seapig-session-saver.service
[Unit]
Description=Scheduler Seapig Session Saver
Requires=scheduler-seapig-server.service
After=scheduler-seapig-server.service
Requires=postgresql.service
After=postgresql.service

[Service]
Environment=RAILS_ENV=production
ExecStart=/usr/bin/bundle exec seapig-session-saver http://127.0.0.1:3001/seapig
Type=simple
Restart=on-failure

WorkingDirectory=/opt/schedy/scheduler
User=schedy
Group=schedy

StartLimitInterval=10s
StartLimitBurst=25
RestartSec=500ms

[Install]
WantedBy=multi-user.target
