# deploy/scheduler-seapig-postgres-notifier.service.in
[Unit]
Description=Scheduler Seapig Postgres Notifier
Requires=scheduler-seapig-server.service
After=scheduler-seapig-server.service

[Service]
Environment=RAILS_ENV=production
ExecStart=/usr/bin/bundle exec seapig-postgresql-notifier
Type=simple
Restart=always

WorkingDirectory=@INSTALL_DIR@
User=schedy
Group=schedy

StartLimitInterval=10s
StartLimitBurst=25
RestartSec=500ms

[Install]
WantedBy=multi-user.target
