# deploy/scheduler-seapig-server.service.in
[Unit]
Description=Scheduler Seapig-Server

[Service]
ExecStart=/usr/bin/bundle exec seapig-server --disable-json-diff
Type=notify
Restart=always

WorkingDirectory=@INSTALL_DIR@
User=schedy
Group=schedy

StartLimitInterval=10s
StartLimitBurst=25
RestartSec=500ms

[Install]
WantedBy=multi-user.target
