# deploy/scheduler-init.service
[Unit]
Description=Scheduler init service
Before=scheduler-seapig-server.service
Requires=postgresql.service
After=postgresql.service

[Service]
Type=oneshot

ExecStart=/usr/libexec/scheduler-init.sh

[Install]
RequiredBy=scheduler-seapig-server.service
